# src/main/centos/etc/crontab
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
*/9 * * * * root /usr/local/bin/deleteJenkinsAgent.sh  >> /var/log/cron 2>&1
